# com.voicetoclaude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voicetoclaude</string>
    <key>ProgramArguments</key>
    <array>
        <string>INSTALL_PATH/venv/bin/python3</string>
        <string>INSTALL_PATH/voice_to_claude.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>INSTALL_PATH</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/voice-to-claude.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voice-to-claude.err</string>
</dict>
</plist>
